# supervisord.conf
[supervisord]
nodaemon=true

[program:jasper]
command=java -Djava.ext.dirs=/jasper/lib:/usr/lib/jvm/java-7-openjdk-amd64/jre/lib/ext -jar /jasper/JavaBridge.jar HTTP:8080